# files/etc/crontabs/root
01 23 * * * /usr/sbin/BMMasterUpdate.sh
02 23 * * * /usr/sbin/DMRIDUpdate.sh
10 23 * * * /usr/sbin/YSFHostsUpdate.sh
