# files/etc/crontabs/root
55 18 * * * /usr/sbin/BMMasterUpdate.sh
55 18 * * * /usr/sbin/DMRIDUpdate.sh
55 18 * * * /usr/sbin/YSFHostsUpdate.sh
